# hw/dv/sv/tl_agent/run
#!/bin/bash

DATE=`date +%Y-%m-%d`

# Simulation output directory
OUT="./out_${DATE}"

mkdir -p ${OUT}

export PRJ_DIR="$(pwd)/../../../../.."

COMPILE_CMD="vcs -file ../../script/vcs.compile.option.f -R \
                 -f ./test.vf -full64 \
                 -l $OUT/compile.log  \
                 -Mdir=$OUT/vcs_simv.csrc \
                 -o $OUT/vcs_simv"

# Compilation
${COMPILE_CMD}
